# packaging/linux/guardswitch.service
[Unit]
Description=GuardSwitch automatic WireGuard switch
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/opt/guardswitch/WgAutoswitch.Service
Restart=on-failure
RestartSec=3
# WICHTIG: PrivateTmp NICHT aktivieren. Die IPC zum User-Tray laeuft ueber einen
# Unix-Socket unter /tmp (/tmp/CoreFxPipe_guardswitch). Ein privates /tmp wuerde
# den Tray dauerhaft vom Service abschneiden ("Service nicht erreichbar").
PrivateTmp=no

[Install]
WantedBy=multi-user.target
